# emr-eAPTs/systemd/emr-eapts-dispense.service
[Unit]
Description=EMR eAPTs Dispense Listener
After=network.target openmrs.service
Wants=openmrs.service

[Service]
Type=simple
User=openmrs
Group=openmrs
WorkingDirectory=/opt/emr-eAPTs
EnvironmentFile=/opt/emr-eAPTs/.env
Environment=PORT=3008
Environment=ROUTE=drugDispense
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/emr-eApts-dispenseCase.log
StandardError=append:/var/log/emr-eApt-dispense-error.log

[Install]
WantedBy=multi-user.target
